# iostat/cron.d/iostat-cron.conf
* * * * * root /etc/zabbix/scripts/iostat/iostat-cron.sh
